# PiTvServer/PiTvServer/pitv-server.service.in
[Unit]
Description=PiTV streaming and recording server
After=multi-user.target

[Service]
Type=exec
ExecStart=@SERVICE_EXECUTABLE_CMD@
Restart=always
RestartSec=10s
KillMode=process
TimeoutSec=120s
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
